# config/systemd/user/remove-old-screenshots.service
[Unit]
Description="Delete old screenshots."

[Service]
Type=simple
ExecStart=/home/mange/.local/bin/remove-old-screenshots
