# platforms/shared/duplicity-backup.timer
[Unit]
Description="Run duplicity at some intervals"

[Timer]
WakeSystem=false
OnStartupSec=30m
OnUnitInactiveSec=1h
Wants=local-fs.target

[Install]
WantedBy=timers.target
