# roles/group-workstation/files/systemd/duplicity-backup.service
[Unit]
Description="Run duplicity, as configured on this machine"

[Service]
Type=simple
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=/home/mange/.config/backup/run
User=mange
Group=mange
